# assets/openrgb-server.service
[Unit]
Description=OpenRGB Server for RGB Keyboard Control
Documentation=https://openrgb.org
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/openrgb --server --noautoconnect
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
